# systemd-20210603/calendar-notifier.service
[Service]
Type=oneshot
ExecStart=%h/bin/calendar-notifier
TimeoutStartSec=20s
